# profile.d/00_locale
#
# Force English as language.
#
if locale -a | grep -q en_GB; then
	export LANG=en_GB.UTF-8
else
	export LANG=en_US.UTF-8
fi
export LC_ALL="$LANG"
