# profile.d/10_autoload
#!/usr/bin/env sh

for file in "$HOME/.functions"/*; do
	name=$(basename "$file")
	echo "$name" | grep -qF . && continue
	eval "$name() {
		unset \"$name\"
		. \"$file\"
		\"$name\" \"\$@\"
	}"
done
