# profile.d/50_gpg
#!/usr/bin/env sh

SSH_AUTH_SOCK="$(gpgconf --list-dirs agent-ssh-socket)"
export SSH_AUTH_SOCK
