# profile.d/80_nix
# Single-user Nix installation.
if [ -f "$HOME/.nix-profile/etc/profile.d/nix.sh" ]; then
	source "$HOME/.nix-profile/etc/profile.d/nix.sh"
fi
# Home manager.
if [ -f "$HOME/.nix-profile/etc/profile.d/hm-session-vars.sh" ]; then
	source "$HOME/.nix-profile/etc/profile.d/hm-session-vars.sh"
fi
