# profile.d/dircolors
#!/usr/bin/env sh

[ -f "$HOME/.dircolors" ] || return

eval "$(dircolors "$HOME/.dircolors")"

