# profile.d/misc
#!/usr/bin/env sh

SHELLCHECK_OPTS="-x"
export SHELLCHECK_OPTS
