# profile.d/rbenv
#!/usr/bin/env sh

[ -d "$HOME/.rbenv" ] || return

export RBENV_ROOT="$HOME/.rbenv"
prepend_to_path "$RBENV_ROOT/bin"
eval "$(rbenv init -)"
